# app/ansible/roles/addons/healthcheck/templates/healthcheck.service.j2
[Unit]
Description=apollo healthcheck
[Service]
Type=oneshot
ExecStart=/usr/bin/curl -fsS --retry 5 -o /dev/null {{ ping_url }}