# app/ansible/roles/addons/healthcheck/templates/healthcheck.timer.j2
[Unit]
Description=apollo healthcheck timer
[Timer]
OnCalendar=*:0/10
Persistent=true
[Install]
WantedBy=timers.target