# roles/addons/gc/templates/garbage-collector.timer.j2
[Unit]
Description=apollo garbage-collection timer
[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true
[Install]
WantedBy=timers.target