# roles/addons/backup/templates/restic-prune.timer.j2
[Unit]
Description=apollo backup pruning timer
[Timer]
OnCalendar=weekly
Persistent=true
[Install]
WantedBy=timers.target