# roles/addons/backup/templates/restic-prune.service.j2
[Unit]
Description=apollo backup pruning
[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic prune
EnvironmentFile={{ app_dir }}/restic-config.cfg