# roles/addons/backup/templates/restic.timer.j2
[Unit]
Description=apollo backup timer
[Timer]
OnCalendar=*:0/30
Persistent=true
[Install]
WantedBy=timers.target